# native-host/com.xport.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xport.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__DAEMON_PATH__</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME_DIR__</string>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>XPORT_LOG_LEVEL</key>
        <string>__LOG_LEVEL__</string>
        <key>XPORT_API_URL</key>
        <string>__API_URL__</string>
        <key>XPORT_INGEST_TOKEN</key>
        <string>__INGEST_TOKEN__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME_DIR__/.xport/daemon-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME_DIR__/.xport/daemon-stderr.log</string>
</dict>
</plist>
