# native-host/com.xport.daemon.service
[Unit]
Description=XPort HTTP daemon
After=network.target

[Service]
ExecStart=__PYTHON_PATH__ __DAEMON_PATH__
Environment=HOME=__HOME_DIR__
Environment=PATH=__PATH__
Environment=XPORT_LOG_LEVEL=__LOG_LEVEL__
Environment=XPORT_API_URL=__API_URL__
Environment=XPORT_INGEST_TOKEN=__INGEST_TOKEN__
Environment=XPORT_TRANSCRIBE_COMMAND=__TRANSCRIBE_COMMAND__
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
